# batterymanager.timer
[Unit]
Description=Battery Manager timer unit

[Timer]
OnBootSec=3s
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
